# scripts/com.huan.cowork-dicas.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.huan.cowork-dicas</string>
  <key>ProgramArguments</key>
  <array>
    <string>/usr/local/bin/node</string>
    <string>/Users/wallace.andrade/Decolar/landing-destinos/scripts/cowork-dicas.mjs</string>
  </array>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>7</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>StandardOutPath</key>
  <string>/tmp/huan-cowork-dicas.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/huan-cowork-dicas.err</string>
  <key>WorkingDirectory</key>
  <string>/Users/wallace.andrade/Decolar/landing-destinos</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/usr/bin:/bin:/opt/homebrew/bin:/Users/wallace.andrade/.nvm/versions/node/v20.19.1/bin</string>
    <key>HOME</key>
    <string>/Users/wallace.andrade</string>
  </dict>
</dict>
</plist>
